# dot_config/systemd/user/easyeffects.service
[Unit]
Description=EasyEffects audio effects daemon (tray)
PartOf=graphical-session.target
After=graphical-session.target dms.service
Requisite=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/easyeffects --service-mode
Restart=on-failure
RestartSec=2

[Install]
WantedBy=graphical-session.target
