# dot_config/systemd/user/wallpaper.timer
[Unit]
Description=Wallpaper timer

[Timer]
OnBootSec=1min
OnUnitActiveSec=5min

[Install]
WantedBy=timers.target
