# dot_config/systemd/user/wallpaper.service
[Unit]
Description=Change wallpaper and theme

[Service]
Type=oneshot
ExecStart=%h/.bin/wallpaper.sh
